# BusinessCard/systemd units/bcarddev-killer.path
[Unit]
Description=bcarddev killer watcher

[Path]
PathModified=/home/ubuntu/bcarddev.tuldok.dev/killer.file

[Install]
WantedBy=multi-user.target
