# BusinessCard/systemd units/bcarddev-watcher.path
[Unit]
Description=bcarddev File Watcher

[Path]
PathModified=/home/ubuntu/bcarddev.tuldok.dev/restarter.file

[Install]
WantedBy=multi-user.target
